# systemd/kubelet.path
[Path]
PathChanged=/var/run/kubelet.lock

[Install]
WantedBy=multi-user.target
